# runit/trinit/finish
#!/bin/bash

echo "Runit stopping trinit"
down .
sleep 5
